# zuul-refresh.service
# /etc/systemd/system/zuul-refresh.service
[Unit]
Description=zuul periodic refresh
Documentation=https://github.com/zenofile/zuul
Requires=zuul-main.service
After=zuulmain.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/zuul --verbose --enforce-sandbox --config /usr/local/etc/zuul/config.yaml --template /usr/local/etc/zuul/template.jinja2 refresh
TimeoutStartSec=120

# Security Hardening - Capabilities
CapabilityBoundingSet=CAP_NET_ADMIN
AmbientCapabilities=CAP_NET_ADMIN
NoNewPrivileges=yes

# Security Hardening - Filesystem Protection
ProtectSystem=full
ProtectHome=yes
PrivateTmp=yes
ReadWritePaths=/run
ConfigurationDirectory=zuul
ConfigurationDirectoryMode=0750

# Security Hardening - Kernel Protection
ProtectKernelModules=yes
ProtectKernelLogs=yes
ProtectControlGroups=yes
ProtectProc=invisible
ProcSubset=pid

# Security Hardening - Device Access
PrivateDevices=yes
DevicePolicy=closed

# Security Hardening - Network
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6 AF_NETLINK

# Security Hardening - System Calls
SystemCallFilter=@system-service
SystemCallFilter=~@privileged @resources @obsolete @debug @mount @swap @reboot @module
SystemCallFilter=socket bind connect sendto recvfrom setsockopt getsockopt
SystemCallErrorNumber=EPERM

# Security Hardening - Namespace Isolation
PrivateUsers=no
PrivateNetwork=no
RestrictNamespaces=yes

# Security Hardening - Memory and Execution
MemoryDenyWriteExecute=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes
LockPersonality=yes
RemoveIPC=yes

# Security Hardening - Miscellaneous
UMask=0077
ProtectClock=yes
KeyringMode=private

# Resource Limits
LimitNOFILE=1024
LimitNPROC=64
TasksMax=128
